# dev-app/template/asset.0690e5bddfbbdac1d92af811114688129d623afb75afad8d69ba4dbbce2b80b4/servertest.service_service.service
[Unit]
Description=Start the project zomboid server script as a service
AssertPathExists=/mnt/default/start-server.sh

[Service]
RemainAfterExit=yes
Restart=always
RestartSec=10
StartLimitInterval=0
TimeoutStartSec=0
ExecStart=/bin/sh -c '/mnt/default/start-server.sh -servername default -adminpassword PasswordXYZ -cachedir=/mnt/default'

[Install]
WantedBy=multi-user.target
