# sheeta-games-dev/assets/dist/units/sheeta.service_ebs-unit.service
[Unit]
Description=Start the project zomboid server script as a service
# AssertPathExists=/start-server.sh

# https://docs.aws.amazon.com/cli/latest/reference/ec2/attach-volume.html

[Service]
RemainAfterExit=yes
Restart=always
RestartSec=10
StartLimitInterval=0
TimeoutStartSec=0
ExecStart=/bin/sh -c 'echo "do ebs"'
# ExecStart=/bin/sh -c '/start-server.sh -servername  -adminpassword  -cachedir='

[Install]
WantedBy=multi-user.target
